# mapi/supervisord.conf
[supervisord]
nodaemon=true

[program:grpc-server]
command=/opt/grpc-example/server/server

[program:rest-proxy]
command=/opt/grpc-example/mapi/rest-proxy